# amlogic/common/rootfs/etc/crontabs/root
0 */12 * * * echo 3 > /proc/sys/vm/drop_caches && rm -rf /tmp/luci*
